# deploy/ttyd.service
[Unit]
Description=ttyd terminal daemon
After=network.target

[Service]
Type=simple
User=terminal
Group=terminal
ExecStart=/usr/local/bin/ttyd --interface 127.0.0.1 --port 7681 --base-path /ttyd --writable --url-arg /bin/bash /opt/web-terminal/scripts/attach-terminal-session.sh
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
