# deploy/web-terminal.service
[Unit]
Description=web-terminal auth app
After=network.target

[Service]
Type=simple
User=terminal
Group=terminal
WorkingDirectory=/opt/web-terminal
EnvironmentFile=/opt/web-terminal/.env
ExecStart=/usr/bin/node /opt/web-terminal/app.js
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/web-terminal
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
